# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:brave-search]
command=npx -y supergateway --stdio "npx -y @modelcontextprotocol/server-brave-search" --port 8003 --baseUrl "http://localhost:8003" --ssePath /sse --messagePath /message --cors --healthEndpoint /healthz
environment=BRAVE_API_KEY=%(ENV_BRAVE_API_KEY)s
stdout_logfile=/app/logs/brave-search.log
stderr_logfile=/app/logs/brave-search.error.log
autorestart=true
startretries=3

[program:emqx]
command=npx -y supergateway --stdio "docker run -i --rm -e EMQX_API_URL=%(ENV_EMQX_API_URL)s -e EMQX_API_KEY=%(ENV_EMQX_API_KEY)s -e EMQX_API_SECRET=%(ENV_EMQX_API_SECRET)s benniuji/emqx-mcp-server" --port 8004 --baseUrl "http://localhost:8004" --ssePath /sse --messagePath /message --cors --healthEndpoint /healthz
stdout_logfile=/app/logs/emqx.log
stderr_logfile=/app/logs/emqx.error.log
autorestart=true
startretries=3

[program:filesystem]
command=npx -y supergateway --stdio "npx -y @modelcontextprotocol/server-filesystem /data" --port 8005 --baseUrl "http://localhost:8005" --ssePath /sse --messagePath /message --cors --healthEndpoint /healthz
stdout_logfile=/app/logs/filesystem.log
stderr_logfile=/app/logs/filesystem.error.log
autorestart=true
startretries=3

[program:git]
command=npx -y supergateway --stdio "uvx mcp-server-git" --port 8006 --baseUrl "http://localhost:8006" --ssePath /sse --messagePath /message --cors --healthEndpoint /healthz
stdout_logfile=/app/logs/git.log
stderr_logfile=/app/logs/git.error.log
autorestart=true
startretries=3

[include]
files = /etc/supervisor/conf.d/*.conf